# autostart/macos/com.splitflap.display.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label -->
    <key>Label</key>
    <string>com.splitflap.display</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>server.py</string>
        <string>8001</string>
    </array>

    <!-- Working Directory - CHANGE THIS to your installation path -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/splitflap</string>

    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- Optional: Add MQTT configuration -->
        <!--
        <key>MQTT_BROKER</key>
        <string>localhost</string>
        <key>MQTT_PORT</key>
        <string>1883</string>
        -->
    </dict>

    <!-- Run at load (boot) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive - restart if crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/tmp/splitflap-stdout.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/tmp/splitflap-stderr.log</string>

    <!-- Session type -->
    <key>SessionCreate</key>
    <true/>
</dict>
</plist>